# systemd/loop-archive@.timer
[Unit]
Description=Periodic run of Loop Archive for %I

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
